# init/claude-mailman@.service
# systemd user template unit for the cli-semaphore mailman daemon.
#
# Install:
#   cp init/claude-mailman@.service ~/.config/systemd/user/
#   loginctl enable-linger alex            # so the user manager runs at boot
#   systemctl --user daemon-reload
#   systemctl --user enable --now claude-mailman@surveyor.service
#
# One instance per recipient agent. The instance name ("%i") is the agent
# name, and is passed through to `claude-msg serve --agent`.
#
# Logs go to the user journal:
#   journalctl --user -u claude-mailman@surveyor.service -f

[Unit]
Description=cli-semaphore mailman for %i
Documentation=https://git.frankenbit.de/frankenbit/cli-semaphore
# Run after the user session is set up; pane-id discovery (#12) writes the
# agents table before any mailman instance starts.
After=default.target
PartOf=default.target

[Service]
Type=simple
ExecStart=/usr/local/bin/claude-msg serve --agent %i

# Restart on crash, but back off so a stuck mailman doesn't hot-loop.
Restart=on-failure
RestartSec=2s

# Watchdog: claude-msg serve pings WATCHDOG=1 every loop iteration when
# WATCHDOG_USEC is exported. systemd kills + restarts the unit if no ping
# arrives within WatchdogSec — covers the case where the loop is wedged
# but the process is still nominally alive (e.g. tmux hung).
WatchdogSec=30s
NotifyAccess=main
Type=notify

# Inherit the operator's PATH and HOME via the user manager.
# Override CLAUDE_MSG_DB here if the install put it somewhere bespoke.
Environment=CLAUDE_MSG_DB=/var/lib/cli-semaphore/messages.db

# Journal capture with per-instance tag so journalctl filtering is clean.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-mailman-%i

# WatchdogSec is not used yet — would require sd_notify(WATCHDOG=1) pings
# from the mailman loop. Restart=on-failure covers the crash case.

[Install]
WantedBy=default.target
