# init/tmux-msg-codex-mailman@.service
# systemd user template unit for the tmux-msg-codex mailman daemon (the OpenAI
# Codex CLI adapter, #248). Sibling of tmux-msg-claude-mailman@; each adapter
# targets its own daemon (the unit prefix is the binary name). Codex is a new
# adapter with no deprecation alias.
#
# A Codex agent uses hook-context delivery (#248 (B), ADR-0009): its messages are
# presented by `tmux-msg-codex hook-context` from the agent's own session, so the
# mailman `serve` short-circuits at startup (notifies ready, exits cleanly — it
# does not paste). The unit is still installed so `register --start-mailman` and
# the health/systemctl surfaces have a real unit to target, and so a future
# paste-capable Codex configuration has somewhere to land.
#
# Install:
#   cp init/tmux-msg-codex-mailman@.service ~/.config/systemd/user/
#   loginctl enable-linger alex            # so the user manager runs at boot
#   systemctl --user daemon-reload
#   systemctl --user enable --now tmux-msg-codex-mailman@surveyor.service
#
# One instance per recipient agent. The instance name ("%i") is the agent
# name, and is passed through to `tmux-msg-codex serve --agent`.
#
# Logs go to the user journal:
#   journalctl --user -u tmux-msg-codex-mailman@surveyor.service -f

[Unit]
Description=tmux-msg mailman for %i (codex adapter)
Documentation=https://git.frankenbit.de/frankenbit/tmux-msg
# Run after the user session is set up; pane-id discovery (#12) writes the
# agents table before any mailman instance starts.
After=default.target
PartOf=default.target

[Service]
Type=notify
ExecStart=/usr/local/bin/tmux-msg-codex serve --agent %i

# Restart on crash, but back off so a stuck mailman doesn't hot-loop.
Restart=on-failure
RestartSec=2s

# Watchdog: tmux-msg-codex serve pings WATCHDOG=1 every loop iteration when
# WATCHDOG_USEC is exported. systemd kills + restarts the unit if no ping
# arrives within WatchdogSec — covers the case where the loop is wedged
# but the process is still nominally alive (e.g. tmux hung).
WatchdogSec=30s
NotifyAccess=main

# Inherit the operator's PATH and HOME via the user manager. The DB location is
# resolved by the binary itself under user-home (XDG: $XDG_DATA_HOME/tmux-msg or
# ~/.local/share/tmux-msg/messages.db) — the same path the CLI uses, since both
# run under the operator's UID (#308). No Environment=CLAUDE_MSG_DB directive: it
# was a system-global-path footgun (#293) and is unnecessary once the default is
# user-home. (CLAUDE_MSG_DB is the substrate-level DB-path env var, shared across
# adapters.) To pin a bespoke DB, add `Environment=CLAUDE_MSG_DB=...` here AND
# ensure callers use the same path (they no longer match the default otherwise).

# Journal capture with per-instance tag so journalctl filtering is clean.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tmux-msg-codex-mailman-%i

[Install]
WantedBy=default.target
