# init/tmux-msg-claude-mailman@.service
# systemd user template unit for the tmux-msg-claude mailman daemon (the Claude
# Code adapter). Renamed from claude-mailman@ per #177; install.sh keeps a
# claude-mailman@ → this symlink for the deprecation cycle (removed v0.11.0).
#
# Install:
#   cp init/tmux-msg-claude-mailman@.service ~/.config/systemd/user/
#   loginctl enable-linger alex            # so the user manager runs at boot
#   systemctl --user daemon-reload
#   systemctl --user enable --now tmux-msg-claude-mailman@surveyor.service
#
# One instance per recipient agent. The instance name ("%i") is the agent
# name, and is passed through to `tmux-msg-claude serve --agent`.
#
# Logs go to the user journal:
#   journalctl --user -u tmux-msg-claude-mailman@surveyor.service -f

[Unit]
Description=tmux-msg mailman for %i
Documentation=https://git.frankenbit.de/frankenbit/tmux-msg
# Run after the user session is set up; pane-id discovery (#12) writes the
# agents table before any mailman instance starts.
After=default.target
PartOf=default.target

# Restart-flood cap (alcatraz-infra#40 belt-and-suspenders against alcatraz-infra#39
# pattern). After 5 starts in 60s, systemd marks the unit `failed` and stops
# restarting until manual intervention (`systemctl --user reset-failed`).
# v0.16.1 made the original visitor-mailman flood class structurally impossible
# (#338 disables the unit on unregister; #340 exits with 0 on agent-not-found
# so Restart=on-failure doesn't trigger), but the cap is still useful insurance
# against future failure modes that exit non-zero in a tight loop.
StartLimitBurst=5
StartLimitInterval=60s

[Service]
Type=simple
ExecStart=/usr/local/bin/tmux-msg-claude serve --agent %i

# Restart on crash, but back off so a stuck mailman doesn't hot-loop.
Restart=on-failure
RestartSec=2s

# Watchdog: tmux-msg-claude serve pings WATCHDOG=1 every loop iteration when
# WATCHDOG_USEC is exported. systemd kills + restarts the unit if no ping
# arrives within WatchdogSec — covers the case where the loop is wedged
# but the process is still nominally alive (e.g. tmux hung).
WatchdogSec=30s
NotifyAccess=main
Type=notify

# Inherit the operator's PATH and HOME via the user manager. The DB location is
# resolved by the binary itself under user-home (XDG: $XDG_DATA_HOME/tmux-msg or
# ~/.local/share/tmux-msg/messages.db) — the same path the CLI uses, since both
# run under the operator's UID (#308). No Environment=CLAUDE_MSG_DB directive: it
# was a system-global-path footgun (#293) and is unnecessary once the default is
# user-home. To pin a bespoke DB, add `Environment=CLAUDE_MSG_DB=...` here AND
# ensure callers use the same path (they no longer match the default otherwise).

# Journal capture with per-instance tag so journalctl filtering is clean.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tmux-msg-claude-mailman-%i

# WatchdogSec is not used yet — would require sd_notify(WATCHDOG=1) pings
# from the mailman loop. Restart=on-failure covers the crash case.

[Install]
WantedBy=default.target
